# systemd/system/checkupdate.timer
[Unit]
Description=Timer that fires the checkupdate service

[Timer]
OnCalendar=*-*-* *:00,30:00

[Install]
WantedBy=timers.target
